# voltronic_nodered_js/rootfs/etc/services.d/voltronic/run
#!/usr/bin/with-contenv bash
set -e
echo "[voltronic] Starting Voltronic JS service..."
exec node /opt/voltronic/voltronic.js
