# voltronic_nodered_js/rootfs/etc/services.d/voltronic/run
#!/usr/bin/with-contenv bash
echo "[voltronic] Disabled (Node-RED handles serial)."
sleep infinity
